# systemd/postflow-watcher.service
[Unit]
Description=UNDLM PostFlow Review Workflow Watcher
After=network.target

[Service]
Type=simple
User=home
WorkingDirectory=/Users/home/Library/CloudStorage/OneDrive-Personnel/WORK/_DEV/rl_postflow
ExecStart=/Users/home/Library/CloudStorage/OneDrive-Personnel/WORK/_DEV/rl_postflow/.venv/bin/python scripts/lucidlink_watcher.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/Users/home/Library/CloudStorage/OneDrive-Personnel/WORK/_DEV/rl_postflow

[Install]
WantedBy=multi-user.target
